# files/rc.local
#!/bin/sh

# Enter custom commands here. By default, they will be executed on startup in parallel with the other services.

# SYSTEM TUNER
if [ -x /bin/system-tuner ]; then
    ( sleep 20; setsid nohup /bin/system-tuner >> /var/log/system-tuner.log 2>&1 ) &
fi

